# etc/systemd/system/smbd.service.d/override.conf
[Service]
BindPaths=/mnt/zssd/nas/torrent:/mnt/zhdd/nas/torrent /mnt/zssd/nas/usenet:/mnt/zhdd/nas/usenet
